# contador.service
[Unit]
Description=user counter service.

[Service]
Type=simple
ExecStart=/bin/contador

[Install]
WantedBy=multi-user.target