# lab6.service
[Unit]
Description=Servicio de usuarios de lab6.

[Service]
Type=simple
ExecStart=/bin/lab6

[Install]
WantedBy=multi-user.target